# s6-initstuff/s6-overlay/s6-rc.d/wireguard/up
wg-quick up client1
#echo UP only one line possible